# deploy/self-host/all-in-one/supervisord.conf
[supervisord]
nodaemon=true
logfile=/data/logs/supervisord.log
pidfile=/run/teamclaw/supervisord.pid
childlogdir=/data/logs

[program:postgres]
command=/usr/local/bin/gosu postgres postgres -D /data/postgres -c listen_addresses=127.0.0.1 -p 5432 -c shared_preload_libraries=pg_net,pg_stat_statements
autostart=true
autorestart=true
priority=10
stdout_logfile=/data/logs/postgres.log
stderr_logfile=/data/logs/postgres.err.log

[program:mqtt-auth]
command=/opt/storage/bin/node /opt/teamclaw/mqtt-auth.mjs
autostart=true
autorestart=true
priority=19
environment=SECRETS_FILE="/data/teamclaw/secrets.env"
stdout_logfile=/data/logs/mqtt-auth.log
stderr_logfile=/data/logs/mqtt-auth.err.log

[program:nanomq]
command=/usr/local/bin/nanomq start --conf /run/teamclaw/nanomq/nanomq.conf
autostart=true
autorestart=true
priority=20
stdout_logfile=/data/logs/nanomq.log
stderr_logfile=/data/logs/nanomq.err.log

[program:gotrue]
command=/bin/sh -c 'set -a; . /run/teamclaw/gotrue/env; set +a; exec /usr/local/bin/gotrue'
autostart=true
autorestart=true
startretries=20
priority=25
stdout_logfile=/data/logs/gotrue.log
stderr_logfile=/data/logs/gotrue.err.log

[program:storage]
command=/bin/sh -c 'set -a; . /run/teamclaw/storage/env; set +a; cd /opt/storage && exec /opt/storage/bin/node dist/start/server.js'
autostart=true
autorestart=true
startretries=20
priority=26
stdout_logfile=/data/logs/storage.log
stderr_logfile=/data/logs/storage.err.log

[program:postgrest]
command=/opt/teamclaw/run-postgrest.sh /run/teamclaw/postgrest/env
autostart=true
autorestart=true
startretries=20
priority=30
stdout_logfile=/data/logs/postgrest.log
stderr_logfile=/data/logs/postgrest.err.log

[program:fc]
command=/bin/sh -c 'set -a; . /run/teamclaw/fc/env; set +a; cd /opt/teamclaw/fc && exec /opt/storage/bin/node dist/server.js'
autostart=true
autorestart=true
startretries=30
priority=40
environment=NODE_ENV="production"
stdout_logfile=/data/logs/fc.log
stderr_logfile=/data/logs/fc.err.log

[program:health]
command=/opt/teamclaw/healthcheck.sh --serve
autostart=true
autorestart=true
priority=80
stdout_logfile=/data/logs/health.log
stderr_logfile=/data/logs/health.err.log

[program:caddy]
command=/usr/bin/caddy run --config /run/teamclaw/caddy/Caddyfile --adapter caddyfile
autostart=true
autorestart=true
priority=90
environment=XDG_DATA_HOME="/data/caddy/data",XDG_CONFIG_HOME="/data/caddy/config"
stdout_logfile=/data/logs/caddy.log
stderr_logfile=/data/logs/caddy.err.log
